# config/systemd/kairos-skill-evolver.service
[Unit]
Description=Skill evolver auto-optimizer
After=network.target
[Service]
ExecStart=/usr/bin/python3 -m ai.skill-evolver
WorkingDirectory=/opt/kairos/ai
Restart=always
RestartSec=5
User=kairos
[Install]
WantedBy=multi-user.target
