# config/systemd/kairos-avionics.service
[Unit]
Description=Avionics bus protocols and telemetry standards
After=network.target
[Service]
ExecStart=/usr/bin/kairos-avionics
Restart=always
RestartSec=5
User=kairos
[Install]
WantedBy=multi-user.target
