# config/systemd/kairos-db.service
[Unit]
Description=SQLite vector database and memory bus
After=network.target
[Service]
ExecStart=/usr/bin/kairos-db
Restart=always
RestartSec=5
User=kairos
[Install]
WantedBy=multi-user.target
