# config/systemd/kairos-quantum.service
[Unit]
Description=Quantum gate emulation and simulator
After=network.target
[Service]
ExecStart=/usr/bin/kairos-quantum
Restart=always
RestartSec=5
User=kairos
[Install]
WantedBy=multi-user.target
