# config/systemd/kairos-vision.service
[Unit]
Description=Real-time vision processing and object detection
After=network.target
[Service]
ExecStart=/usr/bin/kairos-vision
Restart=always
RestartSec=5
User=kairos
[Install]
WantedBy=multi-user.target
