# config/systemd/kairos-context-manager.service
[Unit]
Description=Context manager for sliding window
After=network.target
[Service]
ExecStart=/usr/bin/python3 -m ai.context-manager
WorkingDirectory=/opt/kairos/ai
Restart=always
RestartSec=5
User=kairos
[Install]
WantedBy=multi-user.target
